# systemd/user/kanshi-watcher.path
[Path]
PathChanged=%h/.config/kanshi/

[Install]
WantedBy=graphical-session.target
