# systemd/user/kanshi-watcher.service
[Unit]
Description=Kanshi reloader

[Service]
Type=oneshot
ExecStart=/usr/bin/kanshictl reload

[Install]
WantedBy=graphical-session.target
