# systemd/user/mako-watcher.path
[Path]
PathModified=%h/.config/mako/config

[Install]
WantedBy=graphical-session.target
