# systemd/user/waybar-watcher.path
[Path]
PathChanged=%h/.config/waybar/

[Install]
WantedBy=graphical-session.target
